# s6/mecha-agent/run
#!/bin/bash
exec 2>&1

# Fix ownership of volume-mounted config dirs only (not workspace — too slow for large repos)
chown appuser:appuser /home/appuser/.claude /home/appuser/.codex 2>/dev/null
# Workspace ownership: only fix top-level dir, not recursive
chown appuser:appuser /home/appuser/workspace 2>/dev/null

# Ensure .claude.json has onboarding complete (use node for safe JSON handling)
CLAUDE_JSON="/home/appuser/.claude.json"
s6-setuidgid appuser env HOME=/home/appuser node -e '
  const fs = require("fs");
  const p = process.argv[1];
  let data = {};
  try { data = JSON.parse(fs.readFileSync(p, "utf-8")); } catch {}
  data.hasCompletedOnboarding = true;
  fs.writeFileSync(p, JSON.stringify(data, null, 2) + "\n");
' "$CLAUDE_JSON"

# Check for stop request (dashboard-initiated stop) — bring down the entire container
if [ -f /state/stop-requested ]; then
  rm -f /state/stop-requested
  echo "Stop requested — shutting down container"
  s6-svscanctl -t /run/service 2>/dev/null || kill 1
  exit 0
fi

exec s6-setuidgid appuser env HOME=/home/appuser node /app/agent/entry.js
