# s6/mecha-agent/run
#!/bin/bash
exec 2>&1

# Determine run user: use host UID if provided, otherwise appuser (10001)
RUN_UID="${MECHA_HOST_UID:-10001}"
RUN_GID="${MECHA_HOST_GID:-10001}"
RUN_USER="appuser"

# If host UID differs from appuser, create/update a matching user
if [ "$RUN_UID" != "10001" ]; then
  RUN_USER="mechahost"
  if ! id "$RUN_USER" >/dev/null 2>&1; then
    adduser -D -u "$RUN_UID" -h /home/appuser "$RUN_USER" 2>/dev/null || true
  fi
fi

# Fix ownership of volume-mounted config dirs
chown -R "$RUN_UID:$RUN_GID" /home/appuser/.claude /home/appuser/.codex 2>/dev/null
chown "$RUN_UID:$RUN_GID" /home/appuser/workspace /state /config 2>/dev/null
chown "$RUN_UID:$RUN_GID" /config/bot.yaml /state/credentials.yaml 2>/dev/null

# Ensure .claude.json has onboarding complete (use node for safe JSON handling)
CLAUDE_JSON="/home/appuser/.claude.json"
s6-setuidgid "$RUN_USER" env HOME=/home/appuser node -e '
  const fs = require("fs");
  const p = process.argv[1];
  let data = {};
  try {
    const raw = fs.readFileSync(p, "utf-8");
    data = JSON.parse(raw);
  } catch (e) {
    if (e.code === "ENOENT") {
      // File missing — will create fresh
    } else {
      // Malformed — back up before overwriting
      try { fs.copyFileSync(p, p + ".bak"); } catch {}
      console.error("Warning: .claude.json was malformed, backed up to .claude.json.bak");
    }
  }
  data.hasCompletedOnboarding = true;
  fs.writeFileSync(p, JSON.stringify(data, null, 2) + "\n");
' "$CLAUDE_JSON"

# Check for stop request (dashboard-initiated stop) — bring down the entire container
if [ -f /state/stop-requested ]; then
  rm -f /state/stop-requested
  echo "Stop requested — shutting down container"
  s6-svscanctl -t /run/service 2>/dev/null || kill 1
  exit 0
fi

exec s6-setuidgid "$RUN_USER" env HOME=/home/appuser node /app/agent/entry.js
