# s6/mecha-agent/run
#!/command/execlineb -P
fdmove -c 2 1
s6-setuidgid appuser
export HOME /home/appuser
cd /app
node agent/entry.js
