# s6/tailscale-up/up
#!/command/execlineb -P
foreground {
  /bin/bash -c "
    RETRIES=0
    MAX_RETRIES=3
    DELAY=5

    while [ $RETRIES -lt $MAX_RETRIES ]; do
      if [ -S /var/run/tailscale/tailscaled.sock ]; then
        break
      fi
      RETRIES=$((RETRIES + 1))
      echo \"Waiting for tailscaled... ($RETRIES/$MAX_RETRIES)\"
      sleep $DELAY
    done

    if [ ! -S /var/run/tailscale/tailscaled.sock ]; then
      echo 'WARNING: tailscaled not available, starting in local-only mode'
      exit 0
    fi

    TS_ARGS=\"--hostname=mecha-${MECHA_BOT_NAME}\"

    if [ -n \"$MECHA_TS_AUTH_KEY\" ]; then
      export TS_AUTHKEY=\"$MECHA_TS_AUTH_KEY\"
    fi

    if [ -n \"$MECHA_TS_LOGIN_SERVER\" ]; then
      TS_ARGS=\"$TS_ARGS --login-server=$MECHA_TS_LOGIN_SERVER\"
    fi

    /home/appuser/.local/bin/tailscale up $TS_ARGS --accept-routes || {
      echo 'WARNING: tailscale up failed, continuing in local-only mode'
      exit 0
    }

    echo \"Tailscale connected as mecha-${MECHA_BOT_NAME}\"
  "
}
